# reactionary.service
[Unit]
Description=Reactionary Discord Bot
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/Reactionary
ExecStart=/usr/bin/python3 /path/to/Reactionary/bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
